# deploy/auraflow.service
[Unit]
Description=AuraFlow (Bun + Vite React) production server
Documentation=https://github.com/your-org/auraflow
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=auraflow
Group=auraflow
WorkingDirectory=/opt/auraflow

# `bun --smol` shrinks the baseline heap and disables Bun's JIT-forced-opts,
# trading a little throughput for ~30-50MB lower RSS — important on a 1GB box.
# `start` (from package.json) is just: bun --smol run server/index.ts.
# Bun is installed at /usr/local/bin/bun (a real copy, not a /home symlink)
# because ProtectHome=true below masks /home, which would otherwise cause
# systemd's 203/EXEC when ExecStart points into the per-user install.
#
# Runs the bundled server (server/dist/index.js, a single self-contained
# bun-build output). The bundle inlines all deps, so no node_modules is needed
# at runtime — this is what makes the artifact-only deploy possible.
ExecStart=/usr/local/bin/bun --smol run server/dist/index.js

# Restart policy: recover from crashes but don't thrash if it's wedged.
Restart=always
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=5

# Environment + secrets. server/.env is also loaded by the app itself; this
# file is for anything you want systemd to inject (e.g. override PORT).
EnvironmentFile=-/opt/auraflow/server/.env

# Resource ceilings — keeps the process from eating the whole 1GB box.
MemoryMax=640M
MemoryHigh=512M

# Hardening
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
SystemCallArchitectures=native

# Where logs go (journalctl -u auraflow)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=auraflow

[Install]
WantedBy=multi-user.target
